# deploy/fitnfood.conf
; Supervisor program for the FitnFood Telegram bot.
; Install on the server:
;   sudo cp /var/www/fitnfood/deploy/fitnfood.conf /etc/supervisor/conf.d/fitnfood.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;
; One always-on long-polling process. Lightweight: no web server, no DB server.

[program:fitnfood]
command=/var/www/fitnfood/venv/bin/python -m bot
directory=/var/www/fitnfood
user=root
autostart=true
autorestart=true
startsecs=5
startretries=5
; SIGINT -> our main() catches KeyboardInterrupt and shuts down cleanly.
stopsignal=INT
stopwaitsecs=20
environment=PYTHONUNBUFFERED="1"
; Rotate logs so they never fill the disk.
stdout_logfile=/var/www/fitnfood/bot.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/www/fitnfood/bot-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
